# deployment/tclp-risk-id.service
[Unit]
Description=TCLP Risk ID - Contract Climate Risk Identification Service
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/labs-apps/risk_id
EnvironmentFile=/var/www/labs-apps/risk_id/.env

# Install dependencies and start gunicorn
ExecStartPre=/usr/local/bin/poetry install --no-dev
ExecStart=/usr/local/bin/poetry run gunicorn -c /var/www/labs-apps/risk_id/gunicorn.conf.py tclp.app:app

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tclp-risk-id

[Install]
WantedBy=multi-user.target
